# kanban.service
[Unit]
Description=taskwarrior frontend to display tasks in a kanban board

[Service]
ExecStart=kanban server
